# vim/vimwiki-sync.cron
*/30 * * * * ~/.config/nvim/vimwiki-sync.sh 2>> ~/.config/nvim/vimwiki-sync.log
@reboot ~/.config/nvim/vimwiki-sync.sh 2>> ~/.config/nvim/vimwiki-sync.log
